# .config/systemd/user/wallpaper.timer
[Unit]
Description=Wallpaper timer

[Timer]
OnBootSec=7
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
